# proof-of-concept/openwrt-builder/files/etc/init.d/geth
#!/bin/sh /etc/rc.common
# prevent geth from starting automatically
